# bash/lector.service
[Unit]
Description=Launcher GeminiUNL

[Service]
ExecStart=/home/unl2024/bash/start.sh

[Install]
WantedBy= multi-user.target
